# ramdisk/rootfs/overlay/cv1838/etc/profile.d/multi_arch.sh
#export LD_LIBRARY_PATH=/lib:/lib/arm-linux-gnueabihf:/usr/lib:/usr/lib/arm-linux-gnueabihf
