# ramdisk/rootfs/public/crontabs/musl_riscv64/etc/crontabs/root
*/3 * * * * /bin/ntpdate -p 2 0.pool.ntp.org 2.pool.ntp.org pool.ntp.org > /dev/null 2>&1
