# ramdisk/rootfs/overlay/cv180x_glibc_riscv64/etc/crontabs/root
* */1 * * * /usr/sbin/logrotate /etc/logrotate.conf
